# deploy/zhansai.service
[Unit]
Description=Zhansai Platform (创实部信息化管理平台)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/zhansai/backend
EnvironmentFile=-/var/www/zhansai/backend/.env
Environment=PYTHONUNBUFFERED=1
# 单 worker + 4 线程:内存占用最低,且避免多进程并发写 JSON 的锁竞争
ExecStart=/var/www/zhansai/venv/bin/python -m gunicorn -w 1 --threads 4 -b 127.0.0.1:5000 app:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
